# Server Automation/mount-spot-disk.service
[Unit]
Description=Mount Spot Disks
After=network.target

[Service]
ExecStart=/usr/local/bin/start.sh

[Install]
WantedBy=multi-user.target
